# templates/pxe.service.j2
[Unit]
Description=pxe voice control

[Service]
ExecStart=/usr/bin/python {{ pxe_base_path }}/pxe/pxe.py
User=pxe

[Install]
WantedBy=multi-user.target
